# docker/backup/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/backup/supervisord.log
logfile_maxbytes=50MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/backup
user=backup

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
chown=backup:backup

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

# Backup monitoring service
[program:backup-monitor]
command=/usr/bin/node -r ts-node/register /opt/backup/tools/backup-monitor.ts
directory=/opt/backup/tools
user=backup
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/backup/monitor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
environment=NODE_ENV=production,AUTO_START_BACKUP_MONITOR=true

# Cron service for scheduled backups
[program:cron]
command=/usr/sbin/cron -f
user=root
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/backup/cron.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

# Log rotation service
[program:logrotate]
command=/usr/sbin/logrotate /etc/logrotate.d/backup
user=root
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/var/log/backup/logrotate.log

# Backup health checker (runs periodically)
[program:health-check]
command=/opt/backup/healthcheck.sh
user=backup
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/var/log/backup/healthcheck.log

# Group for all backup services
[group:backup-services]
programs=backup-monitor,cron
priority=999
